# etc/systemd/system/cap-alert-failure.service
[Unit]
Description=Log CAP-Alert service failure

[Service]
Type=oneshot
User=asterisk
Group=asterisk
ExecStart=/usr/lib/cap-alert/log-failure
ProtectSystem=strict
ReadWritePaths=/var/log/cap-alert
